# srcs/influxdb/srcs/supervisord.conf
[supervisord]
nodaemon=true

[program:influxdb]
command=sh setup.sh
autostart=true
autorestart=true